# docker/php/rootfs/etc/services.d/cron-scheduler/finish
#!/bin/sh

echo -n "Debouncing cron-schedule..."
sleep 5
echo " done"
