# docker/php/rootfs/etc/services.d/messenger-async/finish
#!/bin/sh

echo -n "Debouncing messenger-async..."
sleep 5
echo " done"
