# docker/php/rootfs/etc/services.d/messenger-async/run
#!/usr/bin/with-contenv sh

cd /srv
exec s6-setuidgid 1000:1000 ./bin/console messenger:consume async -vvv --limit=100